# roles/timer-jobs/templates/madek-produce-weekly-emails-cron.timer
[Unit]
Description=Madek Produce Weekly Emails Cron Service Timer

[Timer]
OnCalendar={{madek_run_nightly_tasks_on}}
RandomizedDelaySec={{madek_run_nightly_tasks_randomized_delay_max_seconds}}
Persistent=true

[Install]
WantedBy=timers.target
